# dist/tumbleweed-snapshot.timer
[Unit]
Description=tumbleweed snapshot

[Timer]
OnCalendar=hourly
Unit=tumbleweed-snapshot.service

[Install]
WantedBy=timers.target
