# dist/tumbleweed-snapshot.service
# http://www.freedesktop.org/software/systemd/man/systemd.service.html

[Unit]
Description=tumbleweed snapshot

[Service]
Type=oneshot
User=tumbleweed-snapshot
Group=tumbleweed-snapshot
EnvironmentFile=/etc/tumbleweed-snapshot.conf
ExecStart=/usr/bin/tumbleweed-snapshot
